# oncare.service
[Unit]
Description=OnCare Medicine Ordering System Gunicorn daemon
After=network.target postgresql.service redis.service

[Service]
User=oncare
Group=www-data
WorkingDirectory=/path/to/your/project
Environment="PATH=/path/to/your/project/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=medicine_ordering_system.settings_production"
ExecStart=/path/to/your/project/venv/bin/gunicorn \
          --config /path/to/your/project/gunicorn.conf.py \
          --bind unix:/path/to/your/project/oncare.sock \
          medicine_ordering_system.wsgi:application

Restart=always
RestartSec=3

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/path/to/your/project/media /path/to/your/project/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oncare

[Install]
WantedBy=multi-user.target

